# scripts/launchd/com.osxcleaner.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.osxcleaner.agent</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/osxcleaner</string>
        <string>clean</string>
        <string>--safety-level</string>
        <string>4</string>
        <string>--system-caches</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/osxcleaner.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/osxcleaner.error.log</string>

    <key>RunAtLoad</key>
    <false/>

    <key>Nice</key>
    <integer>10</integer>

    <key>LowPriorityIO</key>
    <true/>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
